# etc/systemd/system/delete-monitor-config.service
[Unit]
Description=Delete monitor config

[Service]
ExecStart=rm -f /home/seb/.config/monitors.xml

[Install]
WantedBy=multi-user.target
